# deploy/tests/sandbox/etc/systemd/system/bluetooth.service
[Unit]
Description=Bluetooth Service

[Service]
ExecStart=/usr/lib/bluetooth/bluetoothd

[Install]
WantedBy=multi-user.target
